# deploy/bmo-openclaw.service
[Unit]
Description=BMO OpenClaw Gateway (Telegram Bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/home/prismtek/.local/bin/openclaw gateway start
ExecStop=/home/prismtek/.local/bin/openclaw gateway stop
Restart=on-failure
RestartSec=5
EnvironmentFile=%h/.config/bmo-openclaw.env
# Keep the service running even if the user logs out
RemainAfterExit=yes

[Install]
WantedBy=default.target